# debian/DEBIAN/postinst
#!/bin/bash
# postinst script for raven

set -e

case "$1" in
    configure)
        # Create raven user and group
        if ! getent group raven >/dev/null 2>&1; then
            addgroup --system --quiet raven
        fi
        
        if ! getent passwd raven >/dev/null 2>&1; then
            adduser --system --quiet \
                --home /var/lib/raven \
                --no-create-home \
                --ingroup raven \
                --disabled-password \
                --shell /bin/false \
                --gecos "Raven Network Monitoring" \
                raven
        fi

        # Set proper ownership and permissions
        chown -R raven:raven /var/lib/raven
        chown -R raven:raven /var/log/raven
        chown -R raven:raven /usr/lib/raven
        chown root:raven /etc/raven/config.yaml
        chmod 640 /etc/raven/config.yaml
        chmod 755 /var/lib/raven
        chmod 755 /var/log/raven

        # Create data directory structure
        mkdir -p /var/lib/raven/data
        mkdir -p /var/lib/raven/backups
        chown -R raven:raven /var/lib/raven

        # Reload systemd and enable service
        systemctl daemon-reload
        
        # Enable but don't start the service automatically
        systemctl enable raven.service
        
        echo "Raven has been installed successfully!"
        echo ""
        echo "Files installed:"
        echo "  Configuration:     /etc/raven/config.yaml"
        echo "  Data directory:    /var/lib/raven"
        echo "  Log directory:     /var/log/raven"
        echo "  Documentation:     /usr/share/doc/raven/"
        echo "  Examples:          /usr/share/doc/raven/examples/"
        echo ""
        echo "Quick start:"
        echo "  1. Discover your network:"
        echo "     sudo raven-discover -network 192.168.1.0/24 -output /tmp/discovered.yaml"
        echo ""
        echo "  2. Review and merge the generated config:"
        echo "     # Edit /tmp/discovered.yaml as needed, then:"
        echo "     sudo cp /tmp/discovered.yaml /etc/raven/config.yaml"
        echo ""
        echo "  3. Start Raven:"
        echo "     sudo systemctl start raven"
        echo ""
        echo "Web interface: http://localhost:8000"
        echo ""
        echo "Service management:"
        echo "  Status:    sudo systemctl status raven"
        echo "  Logs:      sudo journalctl -u raven -f"
        echo "  Restart:   sudo systemctl restart raven"
        echo ""
        echo "Documentation:"
        echo "  README:    /usr/share/doc/raven/README.md"
        echo "  Config:    /usr/share/doc/raven/Configuration.md"
        echo "  Examples:  /usr/share/doc/raven/examples/"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
